# workers/userbot/services/loyal-userbot-summary-publish.service
[Unit]
Description=Loyal userbot summary publish once

[Service]
Type=oneshot
WorkingDirectory=/path/to/loyal-app/workers/userbot
EnvironmentFile=/path/to/loyal-app/workers/userbot/.env.local
ExecStart=/path/to/.bun/bin/bun run summary:publish:once
