# workers/userbot/services/loyal-userbot-sync.service
[Unit]
Description=Loyal userbot sync once

[Service]
Type=oneshot
WorkingDirectory=/path/to/loyal-app/workers/userbot
EnvironmentFile=/path/to/loyal-app/workers/userbot/.env.local
ExecStart=/path/to/.bun/bin/bun run sync:once
