# workers/userbot/services/loyal-userbot-summary-publish.timer
[Unit]
Description=Run Loyal userbot summary publish daily at 06:10 UTC

[Timer]
OnCalendar=*-*-* 06:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
